# services/systemd-tmpfiles-setup/up
#!/bin/execlineb -P
systemd-tmpfiles --create --remove --boot --exclude-prefix=/dev
